# services/mariadb/run
#!/usr/bin/with-contenv bash
sudo -u abc /usr/bin/mysqld_safe --skip-syslog
#--binlog-format=MIXED


